# unit-files/canhttp.service
[Unit]
Description=CAN-to-HTTP Daemon
After=can_if.service

[Service]
Type=idle
ExecStart=/home/pi/canhttp.py

[Install]
WantedBy=can_if.service
